# crontab
SHELL=/bin/bash
MAILTO=""
0 12,18,21 * * * /bin/bash /home/duck/hw.sh
