# crontab
SHELL=/bin/bash
MAILTO=""
0 6,12,18 * * * /bin/bash /home/duck/hw.sh
